# samplellama.service
[Unit]
Description=Samplellama - Ollama-to-MCP sampling bridge
Documentation=man:samplellama(1)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/samplellama -mcp-transport http
Restart=on-failure
RestartSec=5

# Hardening
DynamicUser=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
NoNewPrivileges=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes
LockPersonality=yes

[Install]
WantedBy=multi-user.target
